# com.telegram.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.telegram.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vladimir/Work/Cursor site/bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/vladimir/Work/Cursor site/bot/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/vladimir/Work/Cursor site/bot/output.log</string>
</dict>
</plist> 